# restic-sh.timer
[Unit]
Description=Run restic-sh-backup-all every 4 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
